# distribution/com.leaseflex.followup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.leaseflex.followup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/justin/Desktop/leaseflex/distribution/auto_followup.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/justin/Desktop/leaseflex/distribution</string>

    <key>StartInterval</key>
    <integer>21600</integer>

    <key>StandardOutPath</key>
    <string>/Users/justin/Desktop/leaseflex/distribution/logs/launchd_out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/justin/Desktop/leaseflex/distribution/logs/launchd_err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
